# install/docker/supervisord.pool_blocks_merge_nmc.conf
[program:merge.pool_blocks.nmc]
directory=/work/site
command=php artisan merge:pool-blocks-nmc
autostart=false
autorestart=true
killasgroup=true
user=www-data
