# install/docker/supervisord.pool_stats_merge.conf
[program:merge.pool_stats]
directory=/work/site
command=php artisan merge:pool-stats
autostart=false
autorestart=true
killasgroup=true
user=www-data
