# install/docker/supervisord.newblock_notification_btc.conf
[program:newblock_notification.btc]
directory=/work/site
command=php artisan slack:new-block-btc
autostart=false
autorestart=true
killasgroup=true
user=www-data
