# install/docker/supervisord.pool_share_history_merge.conf
[program:merge.pool_share_history]
directory=/work/site
command=php artisan merge:pool-share-history
autostart=false
autorestart=true
killasgroup=true
user=www-data
