# install/docker/supervisord.pool_blocks_merge_btc.conf
[program:merge.pool_blocks.btc]
directory=/work/site
command=php artisan merge:pool-blocks-btc
autostart=false
autorestart=true
killasgroup=true
user=www-data
