# install/docker/supervisord.newblock_notification_nmc.conf
[program:newblock_notification.nmc]
directory=/work/site
command=php artisan slack:new-block-nmc
autostart=false
autorestart=true
killasgroup=true
user=www-data
